# tools/startonboot/ubuntulinux/sceneserver.conf
# sceneserver - start sceneserver server
#
# This task is run on startup to start sceneserver server in a screen environment
# You may access the running servers log output by executing "roomcontrolserversession"

description     "start sceneserver"

start on startup

task
exec sceneserversession 1
